# tools/SalesAgentWA/com.fluxion.salesagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fluxion.salesagent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Developer/CommandLineTools/Library/Frameworks/Python3.framework/Versions/3.9/Resources/Python.app/Contents/MacOS/Python</string>
        <string>/Volumes/MacSSD - Dati/fluxion/tools/SalesAgentWA/agent.py</string>
        <string>send</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Volumes/MacSSD - Dati/fluxion/tools/SalesAgentWA/logs/launchagent.log</string>
    <key>StandardErrorPath</key>
    <string>/Volumes/MacSSD - Dati/fluxion/tools/SalesAgentWA/logs/launchagent_err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Volumes/MacSSD - Dati/fluxion/tools/SalesAgentWA</string>
</dict>
</plist>
